# kernel-usb-workaround.service
# Placed in /etc/systemd/system/kernel-usb-workaround.service
[Unit]
Description=Kernel USB Workaround Service
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/kernel-usb-workaround
User=root
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
