# kernel-usb-workaround.service
# Placed in /etc/systemd/system/kernel-usb-workaround.service
[Unit]
Description=Kernel USB Workaround Service
After=network.target

[Service]
Type=oneshot
ExecStartPre=/bin/sleep 10
ExecStart=/usr/bin/kernel-usb-workaround
User=root
Restart=on-failure

[Install]
WantedBy=multi-user.target
